# scripts/com.daymate.telegram-agent.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.daymate.telegram-agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON_EXE__</string>
    <string>__PROJECT_ROOT__/telegram_agent.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_ROOT__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__PROJECT_ROOT__/logs/telegram-agent.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__PROJECT_ROOT__/logs/telegram-agent.stderr.log</string>
</dict>
</plist>